# overlay/etc/cron.d/faveo-helpdesk
# Faveo Helpdesk cron job to run any scheduled tasks every minute
* * * * * /usr/bin/php /var/www/faveo-helpdesk/artisan schedule:run > /dev/null 2>&1
